# systemctl/phoneHome.timer
[Unit]
Description=timer for heartbeat

[Timer]
Persistent=true
OnUnitActiveSec=300s
OnBootSec=10s

[Install]
WantedBy=timers.target