# systemctl/phoneHome.service
[Unit]
Description=execution script for the heartbeat program

[Service]
Type=oneshot
ExecStart=/usr/bin/env python2